# systemd/p2p-watchdog.service
[Unit]
Description=P2P Wi-Fi Direct - Connection Watchdog
After=p2p-init.service
Requires=p2p-init.service
BindsTo=p2p-init.service

[Service]
Type=simple
EnvironmentFile=/etc/default/video-node
ExecStart=/usr/local/bin/p2p-watchdog.sh
Restart=always
RestartSec=3

StandardOutput=journal
StandardError=journal
SyslogIdentifier=p2p-watchdog

[Install]
WantedBy=multi-user.target
